# etc/infrastructure/supervisor/supervisor.d/crond.ini
[program:crond]
command=/usr/sbin/crond -f
numprocs=1
autostart=false
autorestart=true
priority=300
startretries=200
startsecs=0