# etc/infrastructure/supervisor/supervisor.d/nginx.ini
[program:nginx]
command=/usr/sbin/nginx
numprocs=1
autostart=true
autorestart=true
priority=200
startretries=200
startsecs=0